# deploy/systemd/tearsheet.service
[Unit]
Description=AltBots Tearsheet API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tearsheet-project
ExecStart=/root/tearsheet-project/venv/bin/uvicorn src.api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
